# launchd/com.d0n9x1n.copilot-relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  TEMPLATE — copilot-relay as a per-user launchd agent.

  This file is consumed by install.sh; do NOT symlink it directly.
  install.sh substitutes __HOME__ -> $HOME and writes the rendered plist
  to ~/Library/LaunchAgents/com.d0n9x1n.copilot-relay.plist, then runs
  launchctl bootstrap.

  What the agent does once loaded:
    - Starts `copilot-relay start` on login. install.sh writes
      ~/.copilot-relay/config.yaml with `claudeSetup: false` so the relay
      does not rewrite our symlinked ~/.claude/settings.json.
    - Re-spawns it if it crashes (KeepAlive throttled to one restart per 10s).
    - Listens on http://127.0.0.1:4142 — Claude Code reads this from
      ANTHROPIC_BASE_URL in ~/.claude/settings.json.

  Logs:
    stdout -> ~/Library/Logs/copilot-relay.out.log
    stderr -> ~/Library/Logs/copilot-relay.err.log
    app log -> ~/.copilot-relay/logs/copilot-relay.log

  Manage with:
    launchctl bootstrap   gui/$(id -u) ~/Library/LaunchAgents/com.d0n9x1n.copilot-relay.plist
    launchctl bootout     gui/$(id -u) ~/Library/LaunchAgents/com.d0n9x1n.copilot-relay.plist
    launchctl kickstart -k gui/$(id -u)/com.d0n9x1n.copilot-relay
    launchctl print       gui/$(id -u)/com.d0n9x1n.copilot-relay
    tail -f ~/Library/Logs/copilot-relay.{out,err}.log ~/.copilot-relay/logs/copilot-relay.log

  Notes:
    - launchd plists don't expand $HOME at runtime, so absolute paths are
      required. install.sh templates them per machine.
    - PATH covers Apple Silicon (/opt/homebrew/bin) and Intel
      (/usr/local/bin) Homebrew prefixes so the same plist works on
      both archs.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.d0n9x1n.copilot-relay</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>copilot-relay</string>
        <string>start</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/copilot-relay.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/copilot-relay.err.log</string>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
