# launchd/com.d0n9x1n.copilot-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  TEMPLATE — copilot-api proxy as a per-user launchd agent.

  This file is consumed by install.sh; do NOT symlink it directly.
  install.sh substitutes __HOME__ -> $HOME and writes the rendered plist
  to ~/Library/LaunchAgents/com.d0n9x1n.copilot-api.plist, then runs
  launchctl bootstrap.

  What the agent does once loaded:
    - Starts `copilot-api start --claude-code` on login.
    - Re-spawns it if it crashes (KeepAlive throttled to one restart per 10s).
    - Listens on http://localhost:4141 — Claude Code reads this from
      ANTHROPIC_BASE_URL in ~/.claude/settings.json.

  Logs:
    stdout -> ~/Library/Logs/copilot-api.out.log
    stderr -> ~/Library/Logs/copilot-api.err.log

  Manage with:
    launchctl bootstrap   gui/$(id -u) ~/Library/LaunchAgents/com.d0n9x1n.copilot-api.plist
    launchctl bootout     gui/$(id -u) ~/Library/LaunchAgents/com.d0n9x1n.copilot-api.plist
    launchctl print       gui/$(id -u)/com.d0n9x1n.copilot-api
    tail -f ~/Library/Logs/copilot-api.{out,err}.log

  Notes:
    - launchd plists don't expand $HOME at runtime, so absolute paths are
      required. install.sh templates them per machine.
    - PATH covers Apple Silicon (/opt/homebrew/bin) and Intel (/usr/local/bin)
      Homebrew prefixes so the same plist works on both archs.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.d0n9x1n.copilot-api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/copilot-api</string>
        <string>start</string>
        <string>--claude-code</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/copilot-api.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/copilot-api.err.log</string>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
